# Project/services/ssl-exporter.service
[Unit]
Description=SSL Prometheus Exporter

[Service]
User=ssl-exporter
ExecStart=/usr/local/bin/ssl_exporter

[Install]
WantedBy=multi-user.target
